# web_server_backup.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/local/sbin
MAILTO=someone@example.com
HOME=/home/alextran

0 3 * * * sudo /home/alextran/web-server-backup/backup.sh > /dev/null
